# container/launchd/com.limen.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  SOURCE TEMPLATE — contains NO machine-specific literals (no home-dir paths, no pinned
  python version). Generate the real plist with: scripts/gen-launchd-plist.sh
  All @@PLACEHOLDERS@@ are DERIVED at generation time. Names are outputs, not inputs.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.limen.heartbeat</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/4jp/Workspace/limen/scripts/heartbeat-loop.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key><string>/Users/4jp</string>
        <key>LIMEN_ROOT</key><string>/Users/4jp/Workspace/limen</string>
        <key>LIMEN_WORKDIR</key><string>/Users/4jp/Workspace</string>
        <key>LIMEN_WORKTREES</key><string>/Volumes/Scratch/limen-worktrees</string>
        <key>LIMEN_WORKTREE_ROOT</key><string>/Volumes/Scratch/limen-worktrees</string>
        <key>LIMEN_PYTHON</key><string>/opt/homebrew/Cellar/python@3.14/3.14.6/Frameworks/Python.framework/Versions/3.14/bin/python3.14</string>
        <key>LIMEN_LANES</key><string>codex,opencode,agy,claude,gemini</string>
        <key>LIMEN_DISPATCH_LANES</key><string>auto</string>
        <key>LIMEN_LOCAL_LIMIT</key><string>3</string>
        <key>LIMEN_DISPATCH_ASYNC</key><string>1</string>
        <key>LIMEN_ASYNC_MAX</key><string>10</string>
        <key>LIMEN_VIGILIA</key><string>1</string>
        <key>PATH</key><string>/opt/homebrew/Cellar/python@3.14/3.14.6/Frameworks/Python.framework/Versions/3.14/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>KeepAlive</key><true/>
    <key>RunAtLoad</key><true/>
    <key>ThrottleInterval</key><integer>60</integer>
    <key>ProcessType</key><string>Background</string>
    <key>StandardOutPath</key><string>/Users/4jp/Workspace/limen/logs/heartbeat.out.log</string>
    <key>StandardErrorPath</key><string>/Users/4jp/Workspace/limen/logs/heartbeat.err.log</string>
</dict>
</plist>
