# container/launchd/com.limen.creds-hydrate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  The credential hydration login agent. The metabolize beat already hydrates each daemon cycle;
  this closes the two gaps the beat can't reach: LOGIN time (so a fresh boot has fleet creds before
  any agent runs) and INTERACTIVE/worktree sessions (which aren't the daemon). RunAtLoad fires it at
  login; StartInterval re-runs every 30 min.

  This launchd path is promptless by design: it runs --apply without --op, so it can refresh derived
  lanes (for example, the gh keyring) and already-materialized floors without raising 1Password
  Touch-ID. Static op:// reads require an explicit interactive run with --apply --op, or a promptless
  service-account/Connect setup. Loading is his-hand (one launchctl bootstrap):

    cp container/launchd/com.limen.creds-hydrate.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.limen.creds-hydrate.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.limen.creds-hydrate</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <!-- Mirror the heartbeat's PATH + pinned-venv interpreter so the SAME python3 resolves
             (avoids the macOS TCC python prompt + import drift). `op` lives on the homebrew PATH. -->
        <string>export PATH="/opt/homebrew/bin:/usr/local/bin:$HOME/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin"; PY="$LIMEN_ROOT/.venv/bin/python3"; [ -x "$PY" ] || PY="$(command -v python3 || echo python3)"; exec "$PY" "$LIMEN_ROOT/scripts/creds-hydrate.py" --apply</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key><string>/Users/4jp</string>
        <key>LIMEN_ROOT</key><string>/Users/4jp/Workspace/limen</string>
        <key>PYTHONPATH</key><string>/Users/4jp/Workspace/limen/cli/src</string>
    </dict>
    <!-- RunAtLoad hydrates promptless lanes at login; StartInterval re-runs every 30 min.
         One-shot assess-and-exit model, so no KeepAlive. -->
    <key>StartInterval</key><integer>1800</integer>
    <key>RunAtLoad</key><true/>
    <key>ProcessType</key><string>Background</string>
    <key>StandardOutPath</key><string>/Users/4jp/Workspace/limen/logs/creds-hydrate.out.log</string>
    <key>StandardErrorPath</key><string>/Users/4jp/Workspace/limen/logs/creds-hydrate.err.log</string>
</dict>
</plist>
