# container/launchd/com.limen.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  The dead-man's switch. com.limen.heartbeat is KeepAlive=false (respawns on crash, but NOT if
  it is cleanly killed or wedges mid-beat). This supervisor closes that gap: launchd re-fires it
  every StartInterval REGARDLESS of how the prior run ended, so the watchdog itself can never
  permanently die — and watchdog.py --heal `launchctl kickstart`s the heartbeat whenever it is
  down / not-beating / wedged. heartbeat + watchdog = the flame relights itself with no human hand.

  Self-heal is double-gated by design (watchdog.py): --heal AND LIMEN_WATCHDOG_HEAL=1, both set
  here. Loading is his-hand (one launchctl bootstrap — see docs/NEEDS-HUMAN-DIGEST.md).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.limen.watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <!-- Mirror the heartbeat's PATH + pinned-venv interpreter so watchdog.py resolves the
             SAME python3 the daemon uses (avoids the macOS TCC python prompt + import drift). -->
        <string>export PATH="/opt/homebrew/bin:/usr/local/bin:$HOME/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin"; PY="$LIMEN_ROOT/.venv/bin/python3"; [ -x "$PY" ] || PY="$(command -v python3 || echo python3)"; exec "$PY" "$LIMEN_ROOT/scripts/watchdog.py" --heal</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key><string>/Users/4jp</string>
        <key>LIMEN_ROOT</key><string>/Users/4jp/Workspace/limen</string>
        <key>PYTHONPATH</key><string>/Users/4jp/Workspace/limen/cli/src</string>
        <key>LIMEN_WATCHDOG_HEAL</key><string>1</string>
    </dict>
    <!-- StartInterval re-fires every 5 min no matter what; KeepAlive would fight the one-shot
         assess-and-exit model. RunAtLoad arms it immediately on bootstrap. -->
    <key>StartInterval</key><integer>300</integer>
    <key>RunAtLoad</key><true/>
    <key>ProcessType</key><string>Background</string>
    <key>StandardOutPath</key><string>/Users/4jp/Workspace/limen/logs/watchdog.out.log</string>
    <key>StandardErrorPath</key><string>/Users/4jp/Workspace/limen/logs/watchdog.err.log</string>
</dict>
</plist>
